# deb/DEBIAN/postinst
#!/bin/sh
set -e

systemctl --system daemon-reload >/dev/null || true
systemctl --system enable lovebox >/dev/null || true
systemctl --system start lovebox >/dev/null || true