# deb/DEBIAN/postinst
#!/bin/sh
set -e

systemctl --system daemon-reload >/dev/null || true
systemctl --system enable loveboxpi >/dev/null || true
systemctl --system start loveboxpi >/dev/null || true
