# deb/DEBIAN/prerm
#!/bin/sh
set -e

systemctl --system disable loveboxpi >/dev/null || true
systemctl --system stop loveboxpi >/dev/null || true
